# src/systemd/etl_daemon.service
[Unit]
Description=Medical System ETL Daemon
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=medapp_user
Group=medapp_user
WorkingDirectory=/path/to/your/application
ExecStart=/path/to/your/application/venv/bin/python /path/to/your/application/src/daemons/etl_daemon.py --interval 300 --max-records 1000
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=etl_daemon
Environment=PYTHONUNBUFFERED=1

# Security settings
PrivateTmp=true
ProtectSystem=full
NoNewPrivileges=true
ProtectHome=true

[Install]
WantedBy=multi-user.target