# src/systemd/etl_yottadb_daemon.service
[Unit]
Description=Medical System ETL Daemon - YottaDB (qMS)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=medapp_user
Group=medapp_user
WorkingDirectory=/home/medapp_user/medical_etl
ExecStart=/home/medapp_user/medical_etl/venv/bin/python /home/medapp_user/medical_etl/src/daemons/etl_daemon.py --source yottadb --interval 3600 --max-records 500 --status-file yottadb_status.json
Restart=on-failure
RestartSec=120
TimeoutStartSec=600
TimeoutStopSec=300
StandardOutput=journal
StandardError=journal
SyslogIdentifier=etl_yottadb_daemon
Environment=PYTHONUNBUFFERED=1

# Security settings
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/medapp_user/medical_etl/logs /home/medapp_user/medical_etl/state /home/medapp_user/medical_etl
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target